# deploy/sf-track.service
[Unit]
Description=SF Express tracking proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sf-track
Group=sf-track
WorkingDirectory=/opt/sf-track
EnvironmentFile=/opt/sf-track/.env
ExecStart=/opt/sf-track/bin/sf-track
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/opt/sf-track

[Install]
WantedBy=multi-user.target
